# files/usr/lib/systemd/system/silencer.service
[Unit]
Description=Silence DMESG messages
After=sysinit.target

[Service]
Type=oneshot
ExecStart=dmesg -n 1

[Install]
WantedBy=multiuser.target
