# files/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin f123 --noclear %I $TERM
Type=idle
